# conf/supervisor/conf.d/laalaa.conf
[program:uwsgi-laalaa]
command= /bin/bash /app/docker/start.sh
