# conf/supervisor/conf.d/laalaa.conf
[program:uwsgi-laalaa]
command = /usr/local/bin/uwsgi --ini /etc/uwsgi/laalaa.ini

[program:nginx-laalaa]
command = /usr/sbin/nginx -c /etc/nginx/nginx.conf
